# etc/local.d/plex-reown.start
#!/bin/bash
chown -R plexuser:plexuser /home/plexuser
